# bobbys_table/bobbys_table.service
[Unit]
Description=Run app.py as bobby in a python virtual environment
After=network.target

[Service]
User=bobby
Group=bobby
WorkingDirectory=/home/bobby/SignalWire_python_AI_examples/bobbys_table
ExecStart=/home/bobby/SignalWire_python_AI_examples/bobbys_table/venv/bin/python /home/bobby/SignalWire_python_AI_examples/bobbys_table/app.py
Restart=always
Environment="PATH=/home/bobby/SignalWire_python_AI_examples/bobbys_table/venv/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
